# ccache/.config/profile.d/10-ccache.sh
# shellcheck shell=sh
[ -d /usr/lib/ccache/bin ] && prepend_path /usr/lib/ccache/bin
