# ccache/.config/profile.d/50-ccache.sh
# shellcheck shell=sh
prepend_path /usr/lib/ccache/bin
