# ccache/.config/profile.d/ccache.sh
#!/bin/sh
export CCACHE_DIR=/tmp/ccache
[ -d /usr/lib/ccache/bin ] && PATH=/usr/lib/ccache/bin:$PATH
