# icdiff/.config/profile.d/10-icdiff.sh
# shellcheck shell=sh
export DIFFPROG=icdiff
