# openssl/.config/profile.d/10-openssl.sh
# shellcheck shell=sh
export RANDFILE=${XDG_CACHE_HOME:-$HOME/.cache}/openssl/rnd
