# todotxt/.config/profile.d/10-todotxt.sh
# shellcheck shell=sh
export TODOTXT_CFG_FILE=${XDG_CONFIG_HOME:-$HOME/.config}/todotxt/config
export TODOTXT_DEFAULT_ACTION=ls
