# todotxt/.config/profile.d/todotxt.sh
#!/bin/sh
export TODOTXT_CFG_FILE=${XDG_CONFIG_HOME:-$HOME/.config}/todo/config
export TODOTXT_DEFAULT_ACTION=ls
