# nano/.config/profile.d/10-nano.sh
# shellcheck shell=sh
export EDITOR=nano
