# coreutils/.config/profile.d/10-coreutils.sh
# shellcheck shell=sh
eval "$(dircolors --sh)"
