# wine/.config/profile.d/wine.sh
#!/bin/sh
export WINEPREFIX=${XDG_DATA_HOME:-$HOME/.local/share}/wine
export WINEDLLOVERRIDES=winemenubuilder.exe=d
