# parallel/.config/profile.d/10-parallel.sh
# shellcheck shell=sh
export PARALLEL_HOME=${XDG_CONFIG_HOME:-$HOME/.config}/parallel
