# xdg-utils/.config/profile.d/30-xdg-utils.sh
#!/bin/sh
[ -d ~/.local/bin ] && export PATH=~/.local/bin${PATH:+:$PATH}
