# emacs/.config/profile.d/10-emacs.sh
#!/bin/sh
export VISUAL='emacsclient --tty'
