# emacs/.config/profile.d/10-emacs.sh
# shellcheck shell=sh
export VISUAL='emacsclient --tty'
