# fdroidserver/.config/profile.d/10-fdroidserver.sh
# shellcheck shell=sh
export ANDROID_HOME=/opt/android-sdk
