# gimp/.config/profile.d/gimp.sh
#!/bin/sh
export GIMP2_DIRECTORY=${XDG_CONFIG_HOME:-$HOME/.config}/gimp
