# pacaur/.config/profile.d/pacaur.sh
#!/bin/sh
export AURDEST=${XDG_CACHE_HOME:-$HOME/.cache}/aursync
