# bash/.config/profile.d/90-bash.sh
# shellcheck shell=sh
[ -d ~/.local/bin ] && prepend_path ~/.local/bin
