# bash-completion/.config/profile.d/10-bash-completion.sh
#!/bin/sh
export BASH_COMPLETION_USER_FILE=/usr/share/bash-complete-alias/complete_alias
