# bash-completion/.config/profile.d/bash-completion.sh
#!/bin/sh
export BASH_COMPLETION_USER_FILE=~/.local/lib/complete-alias/bash_completion
