# wget/.config/profile.d/10-wget.sh
# shellcheck shell=sh
export WGETRC=${XDG_CONFIG_HOME:-$HOME/.config}/wget/wgetrc
