# base/.config/profile.d/90-base.sh
# shellcheck shell=sh
prepend_path ~/.local/bin
