# aurutils/.config/profile.d/10-aurutils.sh
# shellcheck shell=sh
export AUR_PAGER=nnn
