# firefox/.config/profile.d/10-firefox.sh
# shellcheck shell=sh
export BROWSER=firefox
